# vtv-scrapy/crontab
*/1 * * * * cd /scrapy_app && scrapy crawl vtvspider >> cron.log 2&>1
